# scripts/bondlink-server.service
[Unit]
Description=Bondlink Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/bondlink-server
Environment="PATH=/opt/bondlink-server/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/opt/bondlink-server/venv/bin/python3 -m server.daemon /etc/bondlink-server/server.yaml
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bondlink-server

# Resource limits
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
